# files/boss.in
#!/bin/sh
#

# PROVIDE: boss
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=boss
rcvar=boss_enable

load_rc_config $name

: ${boss_enable="NO"}
: ${boss_pidfile="/var/run/${name}.pid"}
: ${boss_log="/var/log/boss.log "}

command="/usr/local/bin/${name}"
start_cmd="/usr/sbin/daemon -p ${boss_pidfile} -f -o ${boss_log} ${command} -c /etc/boss.yaml"
pidfile="${boss_pidfile}"

run_rc_command $*
